# tea-api.service
# File path: /etc/systemd/system/tea-api.service
# sudo systemctl daemon-reload
# sudo systemctl start tea-api
# sudo systemctl enable tea-api
# sudo systemctl status tea-api
[Unit]
Description=Tea API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/tea-api  # 注意修改路径
ExecStart=/path/to/tea-api/tea-api --port 3000 --log-dir /path/to/tea-api/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
